# src-webui/port-files/pkg-plist
etc/sudoers.d/appweb-sudoers
share/appweb/dispatcher
share/appweb/fastcgi_params
share/appweb/mime.types
share/appweb/nginx.conf
share/appweb/html/index.php
share/appweb/html/css/_reset.css
share/appweb/html/css/main.css
share/appweb/html/images/pbiicon.php
share/appweb/html/include/Mobile_Detect.php
share/appweb/html/include/header-computer.php
share/appweb/html/include/header-phone.php
share/appweb/html/include/functions.php
share/appweb/html/include/footer-computer.php
share/appweb/html/include/footer-phone.php
share/appweb/html/include/leftmenu.php
share/appweb/html/include/nav-computer.php
share/appweb/html/include/nav-phone.php
share/appweb/html/pages/error.php
share/appweb/html/pages/jails.php
share/appweb/html/pages/main.php
share/appweb/html/pages/notifer.php
share/appweb/html/pages/sysapp.php
share/appweb/html/pages/sysappinfo.php
share/appweb/html/js/doubletaptogo.js
share/appweb/html/js/main.js
etc/rc.d/appweb
@dirrmtry etc/sudoers.d
@dirrm share/appweb/html/css
@dirrm share/appweb/html/images
@dirrm share/appweb/html/include
@dirrm share/appweb/html/js
@dirrm share/appweb/html/pages
@dirrm share/appweb/html
@dirrm share/appweb
