# src-webui/port-files/pkg-plist
etc/sudoers.d/appweb-sudoers
share/appweb/dispatcher
share/appweb/fastcgi_params
share/appweb/mime.types
share/appweb/nginx.conf
share/appweb/html/index.php
share/appweb/html/css/_reset.css
share/appweb/html/css/main.css
share/appweb/html/css/tabs.css
share/appweb/html/images/appcafe.png
share/appweb/html/images/application-exit.png
share/appweb/html/images/categories.png
share/appweb/html/images/configure.png
share/appweb/html/images/dialog-ok.png
share/appweb/html/images/info-tips.png
share/appweb/html/images/install.png
share/appweb/html/images/jail.png
share/appweb/html/images/notifier.png
share/appweb/html/images/pbiicon.php
share/appweb/html/images/rating-0.png
share/appweb/html/images/rating-1.png
share/appweb/html/images/rating-2.png
share/appweb/html/images/rating-3.png
share/appweb/html/images/rating-4.png
share/appweb/html/images/rating-5.png
share/appweb/html/images/remove.png
share/appweb/html/images/restart.png
share/appweb/html/images/start.png
share/appweb/html/images/switch.png
share/appweb/html/images/system.png
share/appweb/html/images/warden.png
share/appweb/html/images/working.gif
share/appweb/html/include/Mobile_Detect.php
share/appweb/html/include/header.php
share/appweb/html/include/functions.php
share/appweb/html/include/footer.php
share/appweb/html/include/globals.php
share/appweb/html/include/leftmenu.php
share/appweb/html/include/nav-computer.php
share/appweb/html/include/nav-phone.php
share/appweb/html/pages/appcafe.php
share/appweb/html/pages/appcafe-browse.php
share/appweb/html/pages/appinfo.php
share/appweb/html/pages/dispatcher.php
share/appweb/html/pages/error.php
share/appweb/html/pages/jails.php
share/appweb/html/pages/main.php
share/appweb/html/pages/notifier.php
share/appweb/html/pages/sysapp.php
share/appweb/html/js/appcafe.js
share/appweb/html/js/doubletaptogo.js
share/appweb/html/js/jquery.ba-hashchange.min.js
share/appweb/html/js/jquery.easytabs.min.js
share/appweb/html/js/jquery.jscroll.min.js
share/appweb/html/js/jquery.min.js
share/appweb/html/js/main.js
etc/rc.d/appweb
@dirrmtry etc/sudoers.d
@dirrm share/appweb/html/css
@dirrm share/appweb/html/images
@dirrm share/appweb/html/include
@dirrm share/appweb/html/js
@dirrm share/appweb/html/pages
@dirrm share/appweb/html
@dirrm share/appweb
