# pbi-modules/ftp/vsftpd/resources/vsftpd
#!/bin/sh

vsftpd_enable=${vsftpd_enable:-"NO"}

. /etc/rc.subr

name=vsftpd
rcvar=`set_rcvar`

prefix=/usr/local
command=${prefix}/libexec/vsftpd

load_rc_config ${name}
run_rc_command "$1" &
