# src-sh/pcbsd-utils/port-files/pkg-plist
bin/about
bin/de-info
bin/de-logout
bin/disable-tor-mode
bin/enable-tor-mode
bin/pc-fbsdupdatecheck
bin/pc-su
bin/pc-krbconf
bin/pc-ldapconf
bin/pc-nssconf
bin/pc-nssldapconf
bin/pc-pamconf
bin/pc-sambaconf
bin/pc-adctl
bin/pc-ldapctl
bin/pc-thinclient
bin/pc-xdgutil
bin/pc-firstgui
bin/pc-runxgui
bin/pc-checkxdisplay
bin/pc-installdialog
bin/pc-sysconfig
bin/pc-systemflag
bin/personacrypt
bin/restamp-grub
bin/warden-to-iocage
share/pcbsd/conf/brand-pcbsd.4th
share/pcbsd/conf/brand-trueos.4th
share/pcbsd/conf/freebsd-update.conf
share/pcbsd/conf/pcbsd-mirrors
share/pcbsd/conf/pcbsd-menu-commands.4th
share/pcbsd/conf/pcbsd-menu.rc
share/pcbsd/conf/pbi-mirror
share/pcbsd/conf/pkg-mirror
share/pcbsd/conf/iso-mirror
share/pcbsd/conf/update-mirror
share/pcbsd/conf/loader.conf.pcbsd
share/pcbsd/conf/loader.conf.trueos
share/pcbsd/conf/logo-pcbsd.4th
share/pcbsd/conf/ntpd.conf
share/pcbsd/conf/rc.conf.pcbsd
share/pcbsd/conf/smb4.conf
share/pcbsd/conf/torrc
share/pcbsd/pc-xdgutil/applications.menu
share/pcbsd/pc-xdgutil/xfce-applications.menu
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-toys.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-more.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-multimedia.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-games-board.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-edu-tools.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-utilities-desktop.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-games-roguelikes.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-utilities-accessibility.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-utilities-pim.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-utilities-file.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-system.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-development-webdevelopment.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-utilities-xutils.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-system-terminal.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-unknown.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-information.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-office.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-internet.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-games-logic.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-main.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-games-strategy.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-edu-miscellaneous.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-settingsmenu.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-edu-mathematics.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-education.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-games-arcade.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-games.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-edu-languages.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-games-card.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-graphics.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-internet-terminal.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-development.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-edu-science.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-development-translation.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-utilities.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-science.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-utilities-peripherals.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-games-kids.directory
share/pcbsd/pc-xdgutil/desktop-dirs/pcbsd-editors.directory
share/pcbsd/pc-thinclient/resources/gdm/gdm-rc
share/pcbsd/pc-thinclient/resources/gdm/gdm-pam
share/pcbsd/pc-thinclient/resources/gdm/custom.conf
share/pcbsd/pc-thinclient/resources/scripts/pc-sysinstall.example
share/pcbsd/pc-thinclient/resources/scripts/wallpaper.png
share/pcbsd/pc-thinclient/resources/scripts/tcsinstall.sh
share/pcbsd/pc-thinclient/resources/scripts/tcslogin.sh
share/pcbsd/pc-thinclient/resources/scripts/cardDetect/XF86Config.compat
share/pcbsd/pc-thinclient/resources/scripts/cardDetect/x_config.sh
share/pcbsd/pc-thinclient/resources/scripts/cardDetect/XF86Config.qemu
share/pcbsd/pc-thinclient/resources/scripts/cardDetect/mouse_config.sh
share/pcbsd/pc-thinclient/resources/scripts/cardDetect/XF86Config.default
share/pcbsd/pc-thinclient/resources/scripts/cardDetect/XF86Config.vmware
share/pcbsd/pc-thinclient/resources/etc/hosts
share/pcbsd/pc-thinclient/resources/etc/fstab
share/pcbsd/pc-thinclient/resources/etc/motd
share/pcbsd/pc-thinclient/resources/etc/ttys
share/pcbsd/pc-thinclient/resources/etc/rc.conf
share/pcbsd/pc-thinclient/resources/etc/gettytab
share/pcbsd/pc-thinclient/resources/root/dot.login
share/pcbsd/pc-thinclient/resources/dhcpd/dhcpd.conf
share/pcbsd/pc-thinclient/resources/boot/beastie.4th
share/pcbsd/pc-thinclient/resources/webui/functions.php
share/pcbsd/pc-thinclient/resources/webui/scripts.php
share/pcbsd/pc-thinclient/resources/webui/clients.php
share/pcbsd/pc-thinclient/resources/webui/config.php.dist
share/pcbsd/pc-thinclient/resources/webui/logincheck.php
share/pcbsd/pc-thinclient/resources/webui/index.php
share/pcbsd/pc-thinclient/resources/webui/login.php
share/pcbsd/de-profiles/lumina.profile
share/pcbsd/de-profiles/lxde.profile
share/pcbsd/de-profiles/kde.profile
share/pcbsd/de-profiles/gnome.profile
share/pcbsd/de-profiles/mate.profile
share/pcbsd/de-profiles/cinnamon.profile
share/pcbsd/de-profiles/xfce.profile
share/pcbsd/desktop-defaults/root/.kde4/share/config/kwinrc
share/pcbsd/desktop-defaults/root/.kde4/share/config/nitrogenrc
share/pcbsd/desktop-defaults/root/.kde4/share/config/kdeglobals
share/pcbsd/desktop-defaults/root/.qt/lipstikstylerc
share/pcbsd/desktop-defaults/root/.bashrc
share/pcbsd/desktop-defaults/root/.profile
share/pcbsd/desktop-defaults/root/.cshrc
share/pcbsd/desktop-defaults/root/.config/qtcurve/stylerc
share/pcbsd/desktop-defaults/root/.config/qtcurve/gtk-icons
share/pcbsd/desktop-defaults/root/.vimrc
share/pcbsd/desktop-defaults/boot/grub/pcbsdfont.pf2
share/pcbsd/desktop-defaults/boot/grub/themes/pcbsd/background.png
share/pcbsd/desktop-defaults/boot/grub/themes/pcbsd/terminal_box_w.png
share/pcbsd/desktop-defaults/boot/grub/themes/pcbsd/dejavu_14.pf2
share/pcbsd/desktop-defaults/boot/grub/themes/pcbsd/terminal_box_e.png
share/pcbsd/desktop-defaults/boot/grub/themes/pcbsd/README
share/pcbsd/desktop-defaults/boot/grub/themes/pcbsd/slider_s.png
share/pcbsd/desktop-defaults/boot/grub/themes/pcbsd/boot_menu_n.png
share/pcbsd/desktop-defaults/boot/grub/themes/pcbsd/terminal_box_c.png
share/pcbsd/desktop-defaults/boot/grub/themes/pcbsd/boot_menu_nw.png
share/pcbsd/desktop-defaults/boot/grub/themes/pcbsd/dejavu_12.pf2
share/pcbsd/desktop-defaults/boot/grub/themes/pcbsd/boot_menu_sw.png
share/pcbsd/desktop-defaults/boot/grub/themes/pcbsd/boot_menu_ne.png
share/pcbsd/desktop-defaults/boot/grub/themes/pcbsd/boot_menu_se.png
share/pcbsd/desktop-defaults/boot/grub/themes/pcbsd/slider_n.png
share/pcbsd/desktop-defaults/boot/grub/themes/pcbsd/boot_menu_s.png
share/pcbsd/desktop-defaults/boot/grub/themes/pcbsd/dejavu_bold_14.pf2
share/pcbsd/desktop-defaults/boot/grub/themes/pcbsd/boot_menu_e.png
share/pcbsd/desktop-defaults/boot/grub/themes/pcbsd/slider_c.png
share/pcbsd/desktop-defaults/boot/grub/themes/pcbsd/boot_menu_w.png
share/pcbsd/desktop-defaults/boot/grub/themes/pcbsd/dejavu_16.pf2
share/pcbsd/desktop-defaults/boot/grub/themes/pcbsd/terminal_box_n.png
share/pcbsd/desktop-defaults/boot/grub/themes/pcbsd/blob_w.png
share/pcbsd/desktop-defaults/boot/grub/themes/pcbsd/boot_menu_c.png
share/pcbsd/desktop-defaults/boot/grub/themes/pcbsd/terminal_box_se.png
share/pcbsd/desktop-defaults/boot/grub/themes/pcbsd/terminal_box_ne.png
share/pcbsd/desktop-defaults/boot/grub/themes/pcbsd/terminal_box_s.png
share/pcbsd/desktop-defaults/boot/grub/themes/pcbsd/terminal_box_sw.png
share/pcbsd/desktop-defaults/boot/grub/themes/pcbsd/terminal_box_nw.png
share/pcbsd/desktop-defaults/boot/grub/themes/pcbsd/theme.txt
share/pcbsd/desktop-defaults/boot/grub/themes/pcbsd/dejavu_10.pf2
share/pcbsd/desktop-defaults/etc/hosts.deniedssh
share/pcbsd/desktop-defaults/etc/crontab
share/pcbsd/desktop-defaults/etc/wall_cmos_clock
share/pcbsd/desktop-defaults/etc/devfs.conf
share/pcbsd/desktop-defaults/etc/devfs.rules
share/pcbsd/desktop-defaults/etc/hosts.allow
share/pcbsd/desktop-defaults/etc/timezone
share/pcbsd/desktop-defaults/etc/blacklist
share/pcbsd/desktop-defaults/usr/local/etc/cups/cups-pdf.conf
share/pcbsd/desktop-defaults/usr/local/etc/cups/cupsd.conf
share/pcbsd/desktop-defaults/usr/local/etc/pkg.conf
share/pcbsd/desktop-defaults/usr/local/etc/fscd.conf
share/pcbsd/desktop-defaults/usr/local/etc/sudoers
share/pcbsd/desktop-defaults/usr/local/etc/beadm.conf
share/pcbsd/desktop-defaults/usr/local/etc/pcbsd.conf
share/pcbsd/desktop-defaults/usr/local/etc/foomatic/filter.conf
share/pcbsd/desktop-defaults/usr/local/etc/screenrc
share/pcbsd/desktop-defaults/usr/local/etc/default/grub
share/pcbsd/desktop-defaults/usr/local/etc/polkit-1/localauthority/50-local.d/system.pkla
share/pcbsd/desktop-defaults/usr/local/etc/volman.conf
share/pcbsd/desktop-defaults/usr/share/skel/.fluxbox/menu
share/pcbsd/desktop-defaults/usr/share/skel/.fluxbox/startup
share/pcbsd/desktop-defaults/usr/share/skel/.fluxbox/overlay
share/pcbsd/desktop-defaults/usr/share/skel/.screenrc
share/pcbsd/desktop-defaults/usr/share/skel/.bashrc
share/pcbsd/desktop-defaults/usr/share/skel/dot.face.icon
share/pcbsd/desktop-defaults/usr/share/skel/dot.login
share/pcbsd/desktop-defaults/usr/share/skel/dot.zshrc
share/pcbsd/desktop-defaults/usr/share/skel/.xprofile
share/pcbsd/desktop-defaults/usr/share/skel/dot.mail_aliases
share/pcbsd/desktop-defaults/usr/share/skel/.vimrc
share/pcbsd/desktop-defaults/usr/share/skel/Documents/.directory
share/pcbsd/desktop-defaults/usr/share/skel/.cache/dconf/user
share/pcbsd/desktop-defaults/usr/share/skel/Pictures/.keep_me
share/pcbsd/desktop-defaults/usr/share/skel/dot.ideskrc
share/pcbsd/desktop-defaults/usr/share/skel/dot.rhosts
share/pcbsd/desktop-defaults/usr/share/skel/dot.Xresources
share/pcbsd/desktop-defaults/usr/share/skel/dot.cshrc
share/pcbsd/desktop-defaults/usr/share/skel/.xscreensaver
share/pcbsd/desktop-defaults/usr/share/skel/dot.login_conf
share/pcbsd/desktop-defaults/usr/share/skel/Videos/.directory
share/pcbsd/desktop-defaults/usr/share/skel/dot.profile
share/pcbsd/desktop-defaults/usr/share/skel/.icons/default/index.theme
share/pcbsd/desktop-defaults/usr/share/skel/dot.shrc
share/pcbsd/desktop-defaults/usr/share/skel/.xsession
share/pcbsd/desktop-defaults/usr/share/skel/dot.loginIcon.png
share/pcbsd/desktop-defaults/usr/share/skel/Downloads/.directory
share/pcbsd/desktop-defaults/usr/share/skel/dot.mailrc
share/pcbsd/desktop-defaults/usr/share/skel/.config/qtcurve/stylerc
share/pcbsd/desktop-defaults/usr/share/skel/.config/qtcurve/gtk-icons
share/pcbsd/desktop-defaults/usr/share/skel/.config/dconf/user
share/pcbsd/desktop-defaults/usr/share/skel/.config/autostart/ibus.desktop
share/pcbsd/desktop-defaults/usr/share/skel/.config/autostart/gnome-fallback-mount-helper.desktop
share/pcbsd/desktop-defaults/usr/share/skel/.config/autostart/gnome-settings-daemon.desktop
share/pcbsd/desktop-defaults/usr/share/skel/.config/qtcurve.gtk-icons
share/pcbsd/desktop-defaults/usr/share/skel/.config/user-dirs.dirs
share/pcbsd/desktop-defaults/usr/share/skel/Music/.directory
share/pcbsd/server-defaults/usr/share/skel/dot.cshrc
share/pcbsd/server-defaults/usr/share/skel/dot.shrc
share/pcbsd/server-defaults/usr/share/skel/dot.profile
share/pcbsd/server-defaults/usr/share/skel/.bashrc
share/pcbsd/server-defaults/usr/share/skel/.vimrc
share/pcbsd/server-defaults/usr/local/etc/default/grub
share/pcbsd/server-defaults/usr/local/etc/pcbsd.conf
share/pcbsd/server-defaults/usr/local/etc/screenrc
share/pcbsd/server-defaults/usr/local/etc/beadm.conf
share/pcbsd/server-defaults/usr/local/etc/pkg.conf
share/pcbsd/server-defaults/root/.vimrc
share/pcbsd/server-defaults/root/.bashrc
share/pcbsd/server-defaults/root/.cshrc
share/pcbsd/server-defaults/root/.profile
share/pcbsd/server-defaults/boot/grub/pcbsdfont.pf2
share/pcbsd/server-defaults/boot/grub/themes/pcbsd/blob_w.png
share/pcbsd/server-defaults/boot/grub/themes/pcbsd/dejavu_bold_14.pf2
share/pcbsd/server-defaults/boot/grub/themes/pcbsd/README
share/pcbsd/server-defaults/boot/grub/themes/pcbsd/boot_menu_n.png
share/pcbsd/server-defaults/boot/grub/themes/pcbsd/terminal_box_n.png
share/pcbsd/server-defaults/boot/grub/themes/pcbsd/dejavu_10.pf2
share/pcbsd/server-defaults/boot/grub/themes/pcbsd/terminal_box_nw.png
share/pcbsd/server-defaults/boot/grub/themes/pcbsd/boot_menu_s.png
share/pcbsd/server-defaults/boot/grub/themes/pcbsd/boot_menu_sw.png
share/pcbsd/server-defaults/boot/grub/themes/pcbsd/terminal_box_s.png
share/pcbsd/server-defaults/boot/grub/themes/pcbsd/boot_menu_nw.png
share/pcbsd/server-defaults/boot/grub/themes/pcbsd/terminal_box_sw.png
share/pcbsd/server-defaults/boot/grub/themes/pcbsd/dejavu_16.pf2
share/pcbsd/server-defaults/boot/grub/themes/pcbsd/slider_c.png
share/pcbsd/server-defaults/boot/grub/themes/pcbsd/terminal_box_ne.png
share/pcbsd/server-defaults/boot/grub/themes/pcbsd/boot_menu_se.png
share/pcbsd/server-defaults/boot/grub/themes/pcbsd/boot_menu_ne.png
share/pcbsd/server-defaults/boot/grub/themes/pcbsd/terminal_box_se.png
share/pcbsd/server-defaults/boot/grub/themes/pcbsd/boot_menu_w.png
share/pcbsd/server-defaults/boot/grub/themes/pcbsd/terminal_box_w.png
share/pcbsd/server-defaults/boot/grub/themes/pcbsd/theme.txt
share/pcbsd/server-defaults/boot/grub/themes/pcbsd/dejavu_12.pf2
share/pcbsd/server-defaults/boot/grub/themes/pcbsd/boot_menu_e.png
share/pcbsd/server-defaults/boot/grub/themes/pcbsd/slider_n.png
share/pcbsd/server-defaults/boot/grub/themes/pcbsd/terminal_box_e.png
share/pcbsd/server-defaults/boot/grub/themes/pcbsd/terminal_box_c.png
share/pcbsd/server-defaults/boot/grub/themes/pcbsd/dejavu_14.pf2
share/pcbsd/server-defaults/boot/grub/themes/pcbsd/boot_menu_c.png
share/pcbsd/server-defaults/boot/grub/themes/pcbsd/slider_s.png
share/pcbsd/server-defaults/boot/grub/themes/pcbsd/background.png
share/pcbsd/scripts/mkzfsdir.sh
share/pcbsd/scripts/reset-firewall
share/pcbsd/scripts/reset-firewall-tor
share/pcbsd/scripts/rmzfsdir.sh
share/pcbsd/scripts/sys-init.sh
share/pcbsd/version
etc/pc-activedirectory.conf
etc/pc-ldap.conf
etc/pkg/fingerprints/pcbsd/trusted/pkg.cdn.pcbsd.org.20131209
etc/pkg/repos/pcbsd.conf.dist
etc/rc.activedirectory
etc/rc.ldap
etc/rc.d/pc-activedirectory
etc/rc.d/pc-kerberos
etc/rc.d/pc-kinit
etc/rc.d/pc-nsswitch
etc/rc.d/pc-pam
etc/rc.d/pc-samba
etc/rc.d/pc-ldap
etc/rc.d/pc-nssldap
etc/rc.d/pc-directory
etc/sudoers.d/pcbsd
etc/fonts/local.conf
man/man1/pc-systemflag.1.gz
man/man8/pc-sysconfig.8.gz
man/man8/personacrypt.8.gz
