# src-qt5/xtrafiles/local/etc/xdg/autostart/pcbsd-updater.desktop
[Desktop Entry]
Version=1.0
Name=TrueOS Updater Tray
Exec=/usr/local/share/pcbsd/xdg-auto/checkupdater.sh
Terminal=false
Type=Application
Icon=pcbsd
Categories=
GenericName=
X-GNOME-Autostart-Phase=Applications
X-MATE-Autostart-Phase=Applications
